# config/systemd/user/updates.timer
[Unit]
Description=Run update check hourly and on boot

[Timer]
OnActiveSec=15min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
